# server/stream-proxy.service
[Unit]
Description=IPTV Stream Proxy
After=network.target

[Service]
Type=notify
ExecStart=/usr/bin/python3 -u /home/youruser/stream-proxy.py
Restart=always
RestartSec=5
WatchdogSec=60
RuntimeMaxSec=86400
User=youruser
LogsDirectory=stream-proxy
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
